# infra/launchd/com.aiassistant.llama-lfm2-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aiassistant.llama-lfm2-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>/Users/byoungyoungla/AI_Project/AI_Assistant/infra/scripts/start-llama-cpp-lfm2-server.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/byoungyoungla/AI_Project/AI_Assistant</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/byoungyoungla</string>
        <key>PATH</key>
        <string>/Users/byoungyoungla/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LLAMA_CPP_MODEL_REPO</key>
        <string>LiquidAI/LFM2-24B-A2B-GGUF:Q4_0</string>
        <key>LLAMA_CPP_PORT</key>
        <string>1242</string>
        <key>LLAMA_CPP_CTX_SIZE</key>
        <string>131072</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/aiassistant-llama-lfm2-server.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/aiassistant-llama-lfm2-server.err.log</string>
</dict>
</plist>